# backend/supervisor.conf
; Only program definitions — placed in /etc/supervisor/conf.d/
; The [supervisord] block lives in the package's /etc/supervisor/supervisord.conf

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:horizon]
command=php /var/www/osdhyan/artisan horizon
directory=/var/www/osdhyan
autostart=true
autorestart=true
user=www-data
stopwaitsecs=3600
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:scheduler]
command=/bin/sh -c "while true; do php /var/www/osdhyan/artisan schedule:run >> /dev/null 2>&1; sleep 60; done"
directory=/var/www/osdhyan
autostart=true
autorestart=true
user=www-data
stdout_logfile=/dev/null
stderr_logfile=/dev/null
